# docker/node-supervisor.conf
[program:node]
command = node /datastringer/wizard.js 
startsecs = 5
stopwaitsecs = 3600
stopasgroup = false
killasgroup = true
